# roles/firewall/files/my_firewall.service
[Unit]
Description=My firewall

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/my_firewall.sh

[Install]
WantedBy=multi-user.target